# scripts/systemd/velora-claude-keepalive.timer
[Unit]
Description=Velora Claude-CLI Keep-Alive (alle 4 Stunden)

[Timer]
OnBootSec=5min
OnUnitActiveSec=4h
# RandomizedDelay verhindert dass Keep-Alive genau zum Briefing-Cron-Zeitpunkt feuert.
RandomizedDelaySec=10min
Persistent=true
Unit=velora-claude-keepalive.service

[Install]
WantedBy=timers.target
