# scripts/systemd/velora-claude-keepalive.service
[Unit]
Description=Velora Claude-CLI Token Keep-Alive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=admin
Group=admin
Environment=HOME=/home/admin
ExecStart=/home/admin/velora/scripts/claude_keepalive.sh
StandardOutput=journal
StandardError=journal
